# systemd/smart_lamp_ui.service
[Unit]
Description=Smart Lamp UI System
After=network.target graphical.target
Wants=graphical.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/smart_lamp
Environment="PATH=/home/pi/smart_lamp/.venv/bin"
Environment="DISPLAY=:0"
Environment="QT_QPA_PLATFORM=linuxfb"
ExecStart=/home/pi/smart_lamp/.venv/bin/python run_ui.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# 资源限制
MemoryMax=512M
CPUQuota=80%

[Install]
WantedBy=graphical.target
